# crontab
# =============================================================================
# Примеры cron-задач для iso-sync
# =============================================================================
# Замените /path/to/ на абсолютный путь к репозиторию.
# Оба скрипта защищены flock — пересечение запусков безопасно.
# При запуске не из TTY прогресс-бар автоматически переключается в компактный
# режим (без \r), поэтому логи в cron-выводе не разрастаются.
# =============================================================================

# -----------------------------------------------------------------------------
# Вариант A (рекомендуется) — раздельное расписание
# -----------------------------------------------------------------------------
# Хэши пересчитываются каждый час: ловим ручные изменения файлов в files/
# (например, если положили ISO туда руками без update_iso.php).
# Обновление образов — раз в сутки в полночь.
# -----------------------------------------------------------------------------

0 * * * *  php /path/to/generate_all_hashes.php > /dev/null 2>&1
0 0 * * *  php /path/to/update_iso.php          > /dev/null 2>&1

# -----------------------------------------------------------------------------
# Вариант B — только обновление (без отдельного крона на хэши)
# -----------------------------------------------------------------------------
# update_iso.php в конце сам запускает generate_all_hashes.php, так что для
# минимальной конфигурации достаточно одной строки. Минус: если файлы
# меняются между запусками update_iso (например, ручное копирование),
# хэши обновятся только при следующем прогоне.
# -----------------------------------------------------------------------------

# 0 0 * * 6  php /path/to/update_iso.php > /dev/null 2>&1   # раз в неделю, ночь субботы
